# mbsync/.local/share/systemd/user/mbsync.timer
[Unit]
Description=Mbsync. Run every minuite.

[Timer]
OnBootSec=1m
OnUnitActiveSec=1m 

[Install]
WantedBy=timers.target
