# .local/share/systemd/user/vdirsyncer.timer
[Unit]
Description=Mbsync. Run every minuite.

[Timer]
OnBootSec=1s
OnUnitActiveSec=1m 

[Install]
WantedBy=timers.target
