# vdirsyncer/.local/share/systemd/user/vdirsyncer.timer
[Unit]
Description=Mbsync. Run every minuite.

[Timer]
OnBootSec=3m
OnUnitActiveSec=1m 

[Install]
WantedBy=timers.target
